# docs/cron/akmods_check.cron
*/15 * * * * $HOME/scripts/check_akmods.sh
